# daemon/launchd/com.tradingbot.local.dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tradingbot.local.dashboard</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/__USER__/Trading/.venv/bin/bot</string>
        <string>dashboard</string>
        <string>--host</string>
        <string>127.0.0.1</string>
        <string>--port</string>
        <string>8765</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/__USER__/Trading</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>60</integer>

    <key>StandardOutPath</key>
    <string>/Users/__USER__/Trading/data/dashboard.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/__USER__/Trading/data/dashboard.stderr.log</string>
</dict>
</plist>
